# deploy/cloudflared-dojojin.service
[Unit]
Description=Cloudflare Tunnel (dojojin.tech)
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
User=kiseki
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/kiseki/.cloudflared/config-host.yml run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
